# apps/node-agent/systemd/node-agent-logrotate.timer
[Unit]
Description=Daily Farm Node Agent log rotation

[Timer]
OnCalendar=*-*-* 02:15:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
