# apps/node-agent/systemd/node-agent-backup-verify.service
[Unit]
Description=Verify Farm Node Agent backups are current
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/node-agent.env
ExecStart=/usr/local/bin/node-agent-verify-backups
