# apps/node-agent/systemd/node-agent-optional-services.service
[Unit]
Description=Farm Node optional services apply

[Service]
Type=oneshot
ExecStart=/usr/local/bin/node-agent-optional-services
